# playbooks/roles/nginx/files/etc/systemd/system/nginx.service.d/extend_file_limit.conf
# extend the limit of open file descriptors up from the default of 1024 on Ubuntu
#
# part of the default jlm2017 config for nginx
# installed via ansible: do not modify here

[Service]
LimitNOFILE=15000
